# supervisor.conf
[program:pygalaxy]
process_name=pygalaxy
command = /opt/pygalaxy/start.sh
autostart = true
autorestart = true
redirect_stderr = true
stdout_logfile = /opt/pygalaxy/pygalaxy.log